# cron/crontab
BASEDIR=/opt/klima
@reboot $BASEDIR/bin/checkrecsensors.sh
*/3 * * * * $BASEDIR/bin/fanctl.py
