# cron/crontab
BASEDIR=/opt/klima
@reboot $BASEDIR/bin/checkrecsensors.py
