# cron/crontab
BASEDIR=/opt/klima
@reboot $BASEDIR/bin/checkrecsensors.sh
*/5 * * * * $BASEDIR/bin/fanctl.py
